# mcelery/celery.conf
[program:celery]
command=/root/venv3/bin/celery -A mcelery worker -B -l info
directory=/root/learn-python
stdout_logfile=/root/learn-python/log/celery.log
autorestart=true
redirect_stderr=true